# osquery/packs/cron_persistence_pack.yml
name: cron_persistence_pack
queries:
  - name: user_crontabs
    query: |
      SELECT * FROM crontab;
    interval: 3600
